# infra/justdoit-task.service
[Unit]
Description=JustDoIt task-service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/justdoit
EnvironmentFile=/opt/justdoit/.env
ExecStart=/usr/bin/java -jar /opt/justdoit/task-service.jar
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
